# scripts/leon-left-brain.service
[Unit]
Description=Leon Left Brain — Awareness + Bridge Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deansabr
WorkingDirectory=/home/deansabr/leon-system
Environment=LEON_BRAIN_ROLE=left
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 main.py --left-brain
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
